# pacman-maintenance@.timer
[Unit]
Description=Perform package manager maintainance daily for %i

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target